# services/kevin_reporter.service
[Unit]
Description=Kevin's Reporter Service
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu
ExecStart=/bin/bash /home/ubuntu/primary_repo/services/run_reporter.sh
Restart=always
KillSignal=SIGINT
TimeoutStopSec=5

[Install]
WantedBy=multi-user.target
