# services/kevin_requestor.service
[Unit]
Description=Kevin's Response/Calc Creator Service
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu
ExecStart=/bin/bash /home/ubuntu/primary_repo/services/run_requestor.sh
Restart=always
KillSignal=SIGINT
TimeoutStopSec=5

[Install]
WantedBy=multi-user.target
